# data/services/add-sd-card-sync.service
[Unit]
Description=Add a MEGA sync for the SD card.

[Service]
Type=simple
ExecStart=/opt/add_sd_card_sync